# apk/aports/rock3b/rock3b-partition-mount-service/rock3b-partition-mount.initd
#!/sbin/openrc-run
# SPDX-License-Identifier: MIT

name="rock3b-partition-mount"
description="Check a hardware device for partitions with labels, and mount them"
command_background="no"

depend() {
  need devfs
  need localmount
  after modules
}

start() {
  ebegin "Mounting additional storage volumes"
  /usr/sbin/rock3b-partition-mount
  eend $?
}
