# overlayfs/inittab
# /etc/inittab - Alpine Linux init configuration
# 
# Root filesystem is mounted read-only by default via overlayfs
# To remount RW for live patching: mount -o remount,rw /

::sysinit:/sbin/openrc boot
::wait:/sbin/openrc default
::ctrlaltdel:/sbin/reboot
::shutdown:/sbin/openrc shutdown
::shutdown:/sbin/halt

tty1::respawn:/sbin/getty -l /bin/sh 38400 tty1
tty2::respawn:/sbin/getty -l /bin/sh 38400 tty2
tty3::respawn:/sbin/getty -l /bin/sh 38400 tty3
tty4::respawn:/sbin/getty -l /bin/sh 38400 tty4
tty5::respawn:/sbin/getty -l /bin/sh 38400 tty5
tty6::respawn:/sbin/getty -l /bin/sh 38400 tty6